# build/rootfs/etc/gearbox/services/postgresql/run
#!/bin/bash
# WARNING: This file has been auto-generated. DO NOT EDIT: WARNING
# Created on 2020-02-17T19:32:10+1100, using template:run.tmpl and json:docker-postgresql.json

# See gearboxworks/gearbox-base for details.
. /build/include-me.sh

c_ok "Gearbox[postgresql]: Starting."

GBEP="$(which ${GEARBOX_ENTRYPOINT})"
if [ "${GBEP}" != "" ]
then
	${GBEP} ${GEARBOX_ARGS}
else
	# Add your service entry here, if no entrypoint defined in the ref container.
	echo ""
fi

