# build/services/postgresql/finish
#!/bin/bash
# WARNING: This file has been auto-generated. DO NOT EDIT: WARNING
# Created on 2020-02-18T20:03:12+1100, using template:finish.tmpl and json:postgresql-10.12.json

# Required to avoid services starting too quickly.
sleep 5

c_ok "Stopped."

