# build/services/postgresql/run
#!/bin/bash
# WARNING: This file has been auto-generated. DO NOT EDIT: WARNING
# Created on 2020-02-18T20:03:12+1100, using template:run.tmpl and json:postgresql-10.12.json

c_ok "Starting."

GBEP="$(which ${GEARBOX_ENTRYPOINT})"
if [ "${GBEP}" != "" ]
then
	${GBEP} ${GEARBOX_ARGS}
else
	# Add your service entry here, if no entrypoint defined in the ref container.
	echo ""
fi

